# infra/templates/systemd/orcazap-app.service.tmpl
[Unit]
Description=OrcaZap FastAPI Application
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=${APP_USER}
Group=${APP_USER}
WorkingDirectory=${APP_DIR}
Environment="PATH=${APP_DIR}/venv/bin"
EnvironmentFile=${APP_ENV_FILE}
ExecStart=${APP_DIR}/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orcazap-app

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target


